# docker/docker.service.d/disable-iptables.conf
# /etc/systemd/system/docker.service.d/disable-iptables.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --iptables=false
